# Reactor/update/dev/reactor-update.service
[Unit]
Description=Update Reactor frontend and backend

[Service]
Type=oneshot
ExecStart=/usr/bin/env bash /home/asus/Reactor/update/update_backend.sh
ExecStart=/usr/bin/env bash /home/asus/Reactor/update/update_frontend.sh